# vmstats.timer
[Unit]
Description=vmstats timer - runs daily at midnight
Documentation=https://github.com/vmstan/vmcrawl

[Timer]
# Run daily at midnight
OnCalendar=daily

# If the system was off at midnight, run shortly after boot
Persistent=true

# Run immediately if the last run was missed
AccuracySec=1min

[Install]
WantedBy=timers.target
